# debian/postrm
#!/bin/bash -e
update-rc.d office defaults
update-rc.d aeroo-docs defaults

#DEBHELPER#
